# deploy/avr-calibration.service
[Unit]
Description=AVR Calibration — web server (Docker)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=pi
ExecStartPre=-/usr/bin/docker rm -f avr-calibration
ExecStart=/usr/bin/docker run --rm \
    --name avr-calibration \
    -p 8000:8000 \
    --device=/dev/bus/usb \
    -v /home/pi/.avr-calibration:/data/.avr-calibration \
    ghcr.io/abarbaccia/avr-calibration:latest
ExecStop=/usr/bin/docker stop avr-calibration
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
